# crontab
0 5 * * * "/home/pibridge/Desktop/Jacob's Ladder Alarm/trigger-alarm.py"
